# system/tupdater.timer
[Unit]
Description=Run tupdater.service every 24 hours

[Timer]
OnCalendar=*:0/24hr